# docker_build/etc/supervisor/conf.d/nubofrontend.conf
[program:nubofrontend]
directory=/opt/nubomanagement-public
command=node dist/multithreadserver-public.js
user=root
autorestart=true
priority=400
