# debbuilder/nubomanagement-public-common/DEBIAN/prerm
#!/bin/sh

if [ "$1" = remove ]; then
  invoke-rc.d nubomanagement-public stop
fi

