# debbuilder/nubomanagement-public-common/DEBIAN/prerm
#!/bin/sh

if [ "$1" = remove ]; then
  systemctl disable nubomanagement-public
fi
systemctl stop nubomanagement-public

